# deploy/klg-eth.service
# Key Level Grid systemd 服务配置 - ETH 实例
# 
# 安装步骤:
# 1. 复制到 systemd 目录:
#    sudo cp /opt/key-level-grid/deploy/klg-eth.service /etc/systemd/system/
#
# 2. 重新加载 systemd:
#    sudo systemctl daemon-reload
#
# 3. 设置开机自启:
#    sudo systemctl enable klg-eth
#
# 4. 启动服务:
#    sudo systemctl start klg-eth
#
# 常用命令:
#    sudo systemctl status klg-eth     # 查看状态
#    sudo systemctl restart klg-eth    # 重启
#    sudo journalctl -u klg-eth -f     # 查看实时日志

[Unit]
Description=Key Level Grid Trading Strategy - ETH
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/key-level-grid
Environment=PATH=/opt/key-level-grid/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/opt/key-level-grid/.env

# 启动命令 - ETH 配置
ExecStart=/opt/key-level-grid/venv/bin/python /opt/key-level-grid/scripts/run/single.py --config /opt/key-level-grid/configs/config_eth.yaml --log-file /opt/key-level-grid/logs/klg_eth.log

# 重启策略
Restart=always
RestartSec=5
TimeoutStartSec=60

# 日志
StandardOutput=append:/opt/key-level-grid/logs/klg_eth_stdout.log
StandardError=append:/opt/key-level-grid/logs/klg_eth_stderr.log
SyslogIdentifier=klg-eth

# 资源限制
MemoryMax=512M
CPUQuota=30%

[Install]
WantedBy=multi-user.target
